# roles/cert-master/templates/certbot-push.service
# ansible-managed
[Unit]
Description=Certbot push to replica (ansible managed)

[Service]
Type=oneshot
ExecStart={{ certbot_push_script }} %i
